# Btalk.service
[Unit]
Description=Btalk

[Service]
ExecStart=/opt/Btalk/bin/Btalk
Restart=always

[Install]
WantedBy=multi-user.target
